# docker/supervisord.conf
; ═══════════════════════════════════════════════════════════════════════════
;  DRE Lab — Supervisord Configuration
;  Manages: PostgreSQL 15 + MongoDB 7 + ttyd (web terminal)
; ═══════════════════════════════════════════════════════════════════════════

[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; ── PostgreSQL 15 ────────────────────────────────────────────────────────
[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main -c config_file=/etc/postgresql/15/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=INT
redirect_stderr=true
stdout_logfile=/var/log/supervisor/postgresql.log
stdout_logfile_maxbytes=10MB

; ── MongoDB 7 ────────────────────────────────────────────────────────────
[program:mongodb]
command=/usr/bin/mongod --config /etc/mongod.conf
user=mongodb
autostart=true
autorestart=true
priority=10
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/supervisor/mongodb.log
stdout_logfile_maxbytes=10MB

; ── ttyd (Web Terminal) ──────────────────────────────────────────────────
[program:ttyd]
command=ttyd --writable --max-clients 5 --port 7681 -t fontSize=14 -t 'theme={"background":"#1e1f29","foreground":"#f8f8f2"}' login -f student
autostart=true
autorestart=true
priority=20
startsecs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/ttyd.log
stdout_logfile_maxbytes=5MB
